# docker/root/etc/services.d/premiumizearr/run
#!/usr/bin/with-contenv bash

cd /opt/app/ || exit

exec \
    s6-setuidgid abc /opt/app/premiumizearrd